# workspace/crontab/laradock
* * * * * /usr/bin/php /var/www/lac/src/artisan schedule:run >> /dev/null 2>&1